# community/modules/logicaldoc-installer/trunk/src/dist/bin/logicaldoc
#! /bin/sh
# Copyright (c) 2008 Logical Objects snc
# All rights reserved.
#
### BEGIN INIT INFO
# Provides:       logicaldoc     
# Required-Start: 
# Should-Start:   
# Required-Stop:
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Description:    LogicalDOC Document Management System
### END INIT INFO

%{INSTALL_PATH}/bin/logicaldoc.sh $*